# systemd/zsh-hist-backup.service
[Unit]
Description=Zsh History Backup

[Service]
Type=oneshot
ExecStart=/home/talati/.local/bin/zsh_hist_backup
