# core/.config/systemd/user/crypt-backup.service
[Unit]
Description=Crypt Backup

[Service]
Type=oneshot
ExecStart=/usr/bin/rclone sync %h/crypt crypt:talati-crypt/crypt
