# core/.config/systemd/user/crypt-backup.timer
[Unit]
Description=Crypt Backup Timer

[Timer]
OnCalendar=hourly
Persistent=true

[Install]
WantedBy=timers.target
