# core/.config/systemd/user/zsh-hist-backup.service
[Unit]
Description=Zsh History Backup

[Service]
Type=oneshot
ExecStart=%h/.local/bin/zsh_hist_backup
